# Pi System Services/hostname.service
[Unit]
Description=Random Hostname Generator
Wants=network-pre.target
After=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /home/admin/script/randhost/hostname.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target