# Pi System Services/unattended-upgrades.timer
[Unit]
Description=Unattended Upgrades Timer

[Timer]
OnBootSec=5min
Unit=unattended-upgrades.service

[Install]
WantedBy=multi-user.target